# deploy/preview-8080.conf
[program:preview-8080]
command=python3.11 web_server.py
directory=/workspace/poc_grader
autostart=true
autorestart=true
startsecs=2
startretries=3
stopsignal=INT
stopwaitsecs=10
stdout_logfile=/tmp/preview-8080.log
stderr_logfile=/tmp/preview-8080.log
redirect_stderr=true
environment=PATH="%(ENV_PATH)s",DASHSCOPE_API_KEY="sk-ws-H.EMDIIYR.jtU9.MEQCIDg63k7FDifjcSOhZIrLlfmhEyb7or87x8Ka3ljuyrKFAiA9kSj93j6TJaUlazt1R_IS1QC-DWan69IoLEyeIbaZhw",BAIDU_API_KEY="6QzUZkERoW31P0kZlpoA8Seh",BAIDU_SECRET_KEY="bmCwZukpPIUxAvssGdS12m9ITj5UhWod",VOLCANO_API_KEY="ark-ddbae8e5-c1ad-4200-8b1d-b8483adca0c6-9eda7"
